# restic-backup@.timer
[Unit]
Description=Backup %i restic repository daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
